# packages/myslowingest/myslowingest.service
[Unit]
Description=MySQL slowlog ingestor

[Service]
Type=forking
User=mysql
Group=mysql
ExecStartPre=/bin/mkdir -p /var/lib/mysql/slowlog/
ExecStart=/usr/bin/myslowingest -d -P /run/myslowingest/myslowingest.pid
Restart=on-failure
RestartSec=30
RuntimeDirectory=myslowingest
PIDFile=/run/myslowingest/myslowingest.pid

[Install]
WantedBy=multi-user.target
